# scripts/ubuntu.upstart
description "storm-mesos"

start on runlevel [2345]
respawn

script
  [ -r /etc/storm/storm-mesos ] && . /etc/storm/storm-mesos
  [ -z "${ULIMIT:-}" ] || ulimit $ULIMIT
  USER="storm"
  exec start-stop-daemon --start -c $USER --name storm-mesos \
    --exec  /usr/local/bin/storm-mesos nimbus >> /var/log/storm/storm-mesos.log 2>&1

end script

